# configs/ssh/supervisor.conf
[program:sshd]
command=sshd -D
autostart=true
autorestart=true
